# profile.d/quartus.sh
# TODO change to quartus path
export PATH=$PATH:/local/altera/14.0/quartus/bin

