# profile.d/vivado.sh
#TODO change to vivado path
export PATH=$PATH:/opt/Xilinx/Vivado/2015.2/bin

